# system-config/templates/rimoo-service2.service.j2
[Unit]
Description=Rimoo Raspberry - Service 2
After=network.target rimoo-app.service
Wants=network.target

[Service]
Type=simple
User={{ system_user }}
Group={{ system_user }}
WorkingDirectory={{ project_path }}/app/service2
ExecStart=/usr/bin/node {{ project_path }}/app/service2/index.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=5
StandardOutput=journal
StandardError=journal

# Variables d'environnement
Environment=NODE_ENV=production
Environment=SERVICE_NAME=service2
Environment=LOG_LEVEL=info

# Limites de ressources
LimitNOFILE=65536
LimitNPROC=4096

# Sécurité
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths={{ project_path }}
ReadWritePaths=/var/log/rimoo
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target 